# deploy/crimson-relay.service
[Unit]
Description=CRIMSON RELAY authoritative game server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
DynamicUser=yes
StateDirectory=crimson-relay
StateDirectoryMode=0700
UMask=0077
WorkingDirectory=/opt/crimson-relay
ExecStart=/usr/bin/node /opt/crimson-relay/src/server/server.js
Environment=NODE_ENV=production
Environment=HOST=127.0.0.1
Environment=PORT=8080
Environment=PERSIST_PATH=/var/lib/crimson-relay/accounts.json
Environment=PERSIST_MANAGE_DIRECTORY=1
EnvironmentFile=-/etc/crimson-relay.env
Restart=on-failure
RestartSec=3s
TimeoutStopSec=45s
NoNewPrivileges=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectSystem=strict
ProtectHome=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
LockPersonality=yes
RestrictRealtime=yes
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
